# deploy/supervisor.conf
[program:app]
command=/app/deploy/launch.sh
directory=/app/config/
user=appuser
environment=HOME=/home/appuser
startsecs=65
